# templates/etc/systemd/system/docker.service.d/environment.conf.j2
# {{ ansible_managed }}

[Service]
{% for env_var in docker__daemon_environment %}
Environment="{{ env_var }}"
{% endfor %}
